# assets/supervisor_mysql.conf
[program:mysql]
command=service mysql start
autostart=true
autorestart=true
